# uvicorn_bd.service
[Unit]
Description=GuvicornAPI for GcS Service

[Service]
Type=simple
WorkingDirectory=/home/ubuntu/FastAPI-BD/
ExecStart=gunicorn main:app -k uvicorn.workers.UvicornWorker
KillMode=mixed
Restart=always

[Install]
WantedBy=multi-user.target